# crontab
30 0 * * * root /overviewer.sh >> /log.txt 2>&1
