# crontab
0 0 * * * root overviewer.py --config=/minecraft/overviewer.cfg && overviewer.py --config=/minecraft/overviewer.cfg --genpoi
